# config-apps/.config/autostart/ghostty.desktop
[Desktop Entry]
Type=Application
Name=Ghostty
Exec=sh -c "tmux has-session -t Development 2>/dev/null || tmuxp load -d development; ghostty -e tmux attach-session -t Development"
Icon=com.mitchellh.ghostty
Comment=Start Ghostty Terminal on Login
X-GNOME-Autostart-enabled=true
X-GNOME-Autostart-Delay=1
