# config-apps/.config/autostart/ghostty.desktop
[Desktop Entry]
Type=Application
Name=Ghostty
Exec=sh -c "(sleep 0.1 && (tmux has-session -t Development 2>/dev/null || tmuxp load -d development &) && sleep 0.7 && ghostty -e tmux attach-session -t Development) &"
Icon=com.mitchellh.ghostty
Comment=Start Ghostty Terminal on Login
X-GNOME-Autostart-enabled=true
