# config-apps/.config/autostart/ghostty.desktop
[Desktop Entry]
Type=Application
Name=Ghostty
Exec=ghostty -e sh -c "i=0; while [ \$\$i -lt 20 ]; do tmux has-session -t Development 2>/dev/null && break; sleep 0.5; i=\$\$((i+1)); done; tmux attach-session -t Development || tmux new-session -s Development"
Icon=com.mitchellh.ghostty
Comment=Start Ghostty Terminal on Login
X-GNOME-Autostart-enabled=false
